# debian/cron.d
# Every 5 minutes, do fpm actions
*/5 * * * *	   root		/usr/lib/alternc/php7-fpm
